# Files/etc/uci-defaults/96-athena-mwan3-ipv6-wan
#!/bin/sh

if ! uci -q show mwan3 >/dev/null 2>&1; then
	exit 0
fi

uci set mwan3.wan6='interface'
uci set mwan3.wan6.enabled='1'
uci set mwan3.wan6.family='ipv6'
uci set mwan3.wan6.reliability='2'
uci -q delete mwan3.wan6.track_ip >/dev/null 2>&1 || true
uci add_list mwan3.wan6.track_ip='2400:3200::1'
uci add_list mwan3.wan6.track_ip='2400:3200:baba::1'
uci add_list mwan3.wan6.track_ip='2606:4700:4700::1111'
uci add_list mwan3.wan6.track_ip='2001:4860:4860::8888'

uci set mwan3.wan6_m1_w3='member'
uci set mwan3.wan6_m1_w3.interface='wan6'
uci set mwan3.wan6_m1_w3.metric='1'
uci set mwan3.wan6_m1_w3.weight='3'

uci set mwan3.wan_only='policy'
uci -q delete mwan3.wan_only.use_member >/dev/null 2>&1 || true
uci add_list mwan3.wan_only.use_member='wan_m1_w3'
uci add_list mwan3.wan_only.use_member='wan6_m1_w3'
uci set mwan3.wan_only.last_resort='unreachable'

uci show mwan3 2>/dev/null | sed -n "s/^\(mwan3\.[^.]*\)=rule$/\1/p" | while read -r section; do
	[ "$section" = "mwan3.default_rule_v6" ] && continue
	[ "$(uci -q get "$section.family")" = "ipv6" ] && continue
	uci set "$section.family=ipv4"
done

uci set mwan3.default_rule_v6='rule'
uci set mwan3.default_rule_v6.dest_ip='::/0'
uci set mwan3.default_rule_v6.family='ipv6'
uci set mwan3.default_rule_v6.use_policy='wan_only'

uci commit mwan3

/etc/init.d/mwan3 enable >/dev/null 2>&1 || true
/etc/init.d/mwan3 restart >/dev/null 2>&1 || true

exit 0
